# stttt2003pk_monitor_agent.ini
[program:run.py]
command=python /data/www/stttt2003pk_monitor_agent_cherry/run.py           
autostart=true            
autorestart=true          
startsecs=10              
startretries=3            
exitcodes=0,2             
stopsignal=QUIT           
stopwaitsecs=10           
user=root              
log_stdout=true           
log_stderr=true           
logfile=/var/log/cherrypy.log  
logfile_maxbytes=1MB      
logfile_backups=10

[program:push_statsd_daemon]
command=python /data/www/stttt2003pk_monitor_agent_cherry/statsd_agent_stttt2003pk.py           
autostart=true            
autorestart=true          
startsecs=10              
startretries=3            
exitcodes=0,2             
stopsignal=QUIT           
stopwaitsecs=10           
user=root              
log_stdout=true           
log_stderr=true           
logfile=/var/log/statsd_agent_stttt2003pk.log  
logfile_maxbytes=1MB      
logfile_backups=10
